# hwagent-api.service
[Unit]
Description=HWAgent API Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/HWAgent
Environment=PATH=/root/HWAgent/.venv/bin
Environment=HWAGENT_VERBOSE=1
ExecStart=/root/HWAgent/.venv/bin/python -m uvicorn api_server:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=3

# Security settings
NoNewPrivileges=true
ReadWritePaths=/root/HWAgent

[Install]
WantedBy=multi-user.target 